# nxshell/module/ocean-display/misc/systemd/ocean-display.service
[Unit]
Description=Lingmo Display service

[Service]
Type=DBus
BusName=org.lingmo.ocean.Display1
ExecStart=/usr/bin/ocean-display
