# nxshell/qtbased/dde-shell/misc/dde-shell@.service.in
[Unit]
Description=Manage an category for dde-shell.

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-shell -C %I

[Install]
WantedBy=multi-user.target
