# init/data_base/openssh-9.8p1-1/files
%FILES%
etc/
etc/pam.d/
etc/pam.d/sshd
etc/ssh/
etc/ssh/moduli
etc/ssh/ssh_config
etc/ssh/ssh_config.d/
etc/ssh/sshd_config
etc/ssh/sshd_config.d/
etc/ssh/sshd_config.d/99-archlinux.conf
usr/
usr/bin/
usr/bin/findssl.sh
usr/bin/scp
usr/bin/sftp
usr/bin/ssh
usr/bin/ssh-add
usr/bin/ssh-agent
usr/bin/ssh-copy-id
usr/bin/ssh-keygen
usr/bin/ssh-keyscan
usr/bin/sshd
usr/lib/
usr/lib/ssh/
usr/lib/ssh/sftp-server
usr/lib/ssh/ssh-keysign
usr/lib/ssh/ssh-pkcs11-helper
usr/lib/ssh/ssh-sk-helper
usr/lib/ssh/sshd-session
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/sshd.service
usr/lib/systemd/system/sshdgenkeys.service
usr/lib/systemd/user/
usr/lib/systemd/user/ssh-agent.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/openssh.conf
usr/share/
usr/share/factory/
usr/share/factory/etc/
usr/share/factory/etc/pam.d/
usr/share/factory/etc/pam.d/sshd
usr/share/factory/etc/ssh/
usr/share/factory/etc/ssh/moduli
usr/share/factory/etc/ssh/ssh_config
usr/share/factory/etc/ssh/sshd_config
usr/share/factory/etc/ssh/sshd_config.d/
usr/share/factory/etc/ssh/sshd_config.d/99-archlinux.conf
usr/share/licenses/
usr/share/licenses/openssh/
usr/share/licenses/openssh/BSD-2-Clause.txt
usr/share/licenses/openssh/LICENCE
usr/share/licenses/openssh/blowfish.BSD-3-Clause.txt
usr/share/licenses/openssh/openbsd-compat.BSD-3-Clause.txt
usr/share/licenses/openssh/openbsd-compat.ISC.txt
usr/share/licenses/openssh/openbsd-compat.MIT.txt
usr/share/licenses/openssh/replacement.BSD-2-Clause.txt
usr/share/licenses/openssh/rijndael.Public-Domain.txt
usr/share/licenses/openssh/snprintf.Public-Domain.txt
usr/share/licenses/openssh/ssh.BSD-3-Clause.txt
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/scp.1.gz
usr/share/man/man1/sftp.1.gz
usr/share/man/man1/ssh-add.1.gz
usr/share/man/man1/ssh-agent.1.gz
usr/share/man/man1/ssh-copy-id.1.gz
usr/share/man/man1/ssh-keygen.1.gz
usr/share/man/man1/ssh-keyscan.1.gz
usr/share/man/man1/ssh.1.gz
usr/share/man/man5/
usr/share/man/man5/moduli.5.gz
usr/share/man/man5/ssh_config.5.gz
usr/share/man/man5/sshd_config.5.gz
usr/share/man/man8/
usr/share/man/man8/sftp-server.8.gz
usr/share/man/man8/ssh-keysign.8.gz
usr/share/man/man8/ssh-pkcs11-helper.8.gz
usr/share/man/man8/ssh-sk-helper.8.gz
usr/share/man/man8/sshd.8.gz
var/
var/empty/

%BACKUP%
etc/pam.d/sshd	6b6a43083dc71f1cd7a16c92a352e792
etc/ssh/ssh_config	9bf7800f4e812a4f047e7f1b4425cf5a
etc/ssh/sshd_config	4bd4603152da5f703d564c9adc418fc0

