# nxshell/module/ocean-session/tools/ocean-keyring-checker/ocean-keyring-checker.service
[Unit]
Description=ocean-keyring-checker service
CollectMode=inactive-or-failed

Requisite=ocean-session-pre.target
After=ocean-session-pre.target

Requisite=ocean-session-initialized.target
PartOf=ocean-session-initialized.target
Before=ocean-session-initialized.target

[Service]
Type=simple
ExecStart=/usr/bin/ocean-keyring-checker
