# init/data_base/iwd-2.19-1/files
%FILES%
usr/
usr/bin/
usr/bin/hwsim
usr/bin/iwctl
usr/bin/iwmon
usr/lib/
usr/lib/iwd/
usr/lib/iwd/ead
usr/lib/iwd/iwd
usr/lib/modules-load.d/
usr/lib/modules-load.d/pkcs8.conf
usr/lib/systemd/
usr/lib/systemd/network/
usr/lib/systemd/network/80-iwd.link
usr/lib/systemd/system/
usr/lib/systemd/system/ead.service
usr/lib/systemd/system/iwd.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/net.connman.ead.service
usr/share/dbus-1/system-services/net.connman.iwd.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/ead-dbus.conf
usr/share/dbus-1/system.d/hwsim-dbus.conf
usr/share/dbus-1/system.d/iwd-dbus.conf
usr/share/doc/
usr/share/doc/iwd/
usr/share/doc/iwd/README
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hwsim.1.gz
usr/share/man/man1/iwctl.1.gz
usr/share/man/man1/iwmon.1.gz
usr/share/man/man5/
usr/share/man/man5/iwd.ap.5.gz
usr/share/man/man5/iwd.config.5.gz
usr/share/man/man5/iwd.network.5.gz
usr/share/man/man7/
usr/share/man/man7/iwd.debug.7.gz
usr/share/man/man8/
usr/share/man/man8/ead.8.gz
usr/share/man/man8/iwd.8.gz

