# nxshell/module/ocean-appearance/misc/systemd/ocean-fakewm.service.in
[Unit]
Description=ocean window manager

Requisite=plasma-kglobalaccel.service

[Service]
Type=simple
ExecCondition=/bin/sh -c 'test "$XDG_SESSION_TYPE" != "wayland" || exit 2'
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ocean-fakewm
Slice=session.slice
