# shell/frameworks/lingmo-workspace/ksmserver/lingmo-ksmserver.service.in
[Unit]
Description=KDE Session Management Server
PartOf=graphical-session.target
After=lingmo-kwin_wayland.service lingmo-kcminit.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ksmserver
BusName=org.kde.ksmserver
Slice=session.slice

[Install]
WantedBy=lingmo-core.target
