# init/data_base/rauc-1.11.3-1/files
%FILES%
etc/
etc/rauc/
usr/
usr/bin/
usr/bin/rauc
usr/lib/
usr/lib/rauc-service.sh
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/rauc.service
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/rauc.conf
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/de.pengutronix.rauc.Installer.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/de.pengutronix.rauc.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/de.pengutronix.rauc.conf
usr/share/doc/
usr/share/doc/rauc/
usr/share/doc/rauc/grub.conf
usr/share/doc/rauc/uboot.sh
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/rauc.1.gz

