# nxshell/qtbased/dde-session/systemd/dde-session-shutdown.service.in
[Unit]
Description=dde-session-shutdown service, trigger shutdown of a deepin Session
CollectMode=inactive-or-failed

PartOf=dde-session.target

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/dde-session-ctl --shutdown
