# nxshell/qtbased/dde-session/systemd/dde-session-manager.target
[Unit]
Description=DDE Session Manager is ready
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requisite=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target

