# nxshell/qtbased/dde-session/systemd/dde-session-initialized.target
[Unit]
Description=dde-session-initialized target, task to initialize the desktop environment.
OnFailure=dde-session-shutdown.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requires=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session.target
PartOf=dde-session.target
Before=dde-session.target
