# nxshell/qtbased/dde-session/systemd/dde-session-pre.target
[Unit]
Description=dde-session-pre target, tasks to be run before the deepin desktop environment is initialized.
OnFailure=dde-session-shutdown.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requires=graphical-session-pre.target
After=graphical-session-pre.target

Requires=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target
