# nxshell/qtbased/dde-session/systemd/dde-session-initialized.target.wants/dde-osd.target
[Unit]
Description=dde-osd target
DefaultDependencies=no
RefuseManualStart=no
RefuseManualStop=no

BindsTo=dde-osd.service
After=dde-osd.service