# distribution/current/lingmo-config/common/includes.chroot_befor_packages/etc/profile.d/qt5-accessibility.sh
export QT_LINUX_ACCESSIBILITY_ALWAYS_ON=1
