# nxshell/qtbased/dde-widgets/app/services/systemd/dde-widgets.service.in
[Unit]
Description=dde-widgets service

Wants=dde-osd.target
After=dde-osd.target

[Service]
Type=simple
RemainAfterExit=no
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-widgets
