# nxshell/module/ocean-application-wizard/systemd/org.lingmo.ocean.daemon.Launcher1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=Legacy OCEAN Application Manager launcher daemon compat background service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/ocean-application-wizard-daemon-compat
Type=dbus
BusName=org.lingmo.ocean.daemon.Launcher1
Restart=on-failure
RestartSec=1s
