# shell/systemui-desktop/lingmo-appstore/data/autostart/org.lingmo.Appstore.desktop.in
[Desktop Entry]
Type=Application
Name=LINGMO Appstore
Exec=@bindir@/lingmo-appstore --gapplication-service
OnlyShowIn=LINGMO;Unity;
NotShowIn=Budgie
NoDisplay=true
