# core32/shell/lingmo-startshell.service.in
[Unit]
Description=Lingmo OS Workspace
After=lingmo-ksmserver.service lingmo-kcminit.service
PartOf=graphical-session.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
ExecStart=lingmo-session
ExecStartPost=/usr/bin/check_start.sh
Restart=no
Slice=session.slice
TimeoutSec=infinity

[Install]
WantedBy=lingmo-core.target