# nxshell/module/ocean-session/systemd/ocean-session-initialized.target
[Unit]
Description=ocean-session-initialized target, task to initialize the desktop environment.
OnFailure=ocean-session-shutdown.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requires=ocean-session-pre.target
After=ocean-session-pre.target

Requisite=ocean-session.target
PartOf=ocean-session.target
Before=ocean-session.target
