# nxshell/module/ocean-session/systemd/ocean-session-shutdown.service.in
[Unit]
Description=ocean-session-shutdown service, trigger shutdown of a lingmo Session
CollectMode=inactive-or-failed

PartOf=ocean-session.target

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/ocean-session-ctl --shutdown
