# nxshell/module/ocean-session/systemd/ocean-session-manager.target
[Unit]
Description=OCEAN Session Manager is ready
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requisite=ocean-session-pre.target
After=ocean-session-pre.target

Requisite=ocean-session-initialized.target
PartOf=ocean-session-initialized.target
Before=ocean-session-initialized.target

