# nxshell/module/ocean-session/systemd/ocean-session.target
[Unit]
Description=ocean-session target
DefaultDependencies=no

BindsTo=graphical-session.target
Before=graphical-session.target

Requires=ocean-session-initialized.target
After=ocean-session-initialized.target

Requires=ocean-session-manager.service
After=ocean-session-manager.service
