# nxshell/module/ocean-session/systemd/ocean-session-pre.target
[Unit]
Description=ocean-session-pre target, tasks to be run before the lingmo desktop environment is initialized.
OnFailure=ocean-session-shutdown.target
OnFailureJobMode=replace-irreversibly
DefaultDependencies=no
RefuseManualStart=yes
RefuseManualStop=yes

Requires=graphical-session-pre.target
After=graphical-session-pre.target

Requires=ocean-session-initialized.target
PartOf=ocean-session-initialized.target
Before=ocean-session-initialized.target
