# nxshell/qtbased/dde-calendar/calendar-service/assets/data/com.dde.calendarserver.calendar.timer
[Unit]
Description = dde-calendar date service

[Timer]
Unit = com.dde.calendarserver.calendar.service
OnActiveSec = 1s
OnUnitInactiveSec = 10m
AccuracySec = 1us
RandomizedDelaySec = 0
OnBootSec = 1m

[Install]
WantedBy = timers.target