# nxshell/qtbased/dde-calendar/calendar-service/assets/data/dde-calendar.service.in
[Unit]
Description = dde-calendar date service

[Service]
ExecStart=@EXEC_DDE_CALENDAR@
Type=dbus
BusName=com.deepin.Calendar
Slice=app.slice