# init/data_base/irqbalance-1.9.4-1/files
%FILES%
etc/
etc/irqbalance.env
usr/
usr/bin/
usr/bin/irqbalance
usr/bin/irqbalance-ui
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/irqbalance.service
usr/share/
usr/share/doc/
usr/share/doc/irqbalance/
usr/share/doc/irqbalance/AUTHORS
usr/share/doc/irqbalance/README.md
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/irqbalance-ui.1.gz
usr/share/man/man1/irqbalance.1.gz

%BACKUP%
etc/irqbalance.env	fd441b2b501b7258a1043ba135736623

