# init/data_base/dnf-4.21.0-1/files
%FILES%
etc/
etc/dnf/
etc/dnf/aliases.d/
etc/dnf/aliases.d/zypper.conf
etc/dnf/automatic.conf
etc/dnf/dnf.conf
etc/dnf/protected.d/
etc/dnf/protected.d/yum.conf
etc/logrotate.d/
etc/logrotate.d/dnf
usr/
usr/bin/
usr/bin/dnf
usr/bin/dnf-3
usr/bin/dnf-automatic
usr/bin/dnf-automatic-3
usr/lib/
usr/lib/python3.12/
usr/lib/python3.12/site-packages/
usr/lib/python3.12/site-packages/dnf-4.21.0.dist-info/
usr/lib/python3.12/site-packages/dnf-4.21.0.dist-info/METADATA
usr/lib/python3.12/site-packages/dnf/
usr/lib/python3.12/site-packages/dnf/__init__.py
usr/lib/python3.12/site-packages/dnf/automatic/
usr/lib/python3.12/site-packages/dnf/automatic/__init__.py
usr/lib/python3.12/site-packages/dnf/automatic/emitter.py
usr/lib/python3.12/site-packages/dnf/automatic/main.py
usr/lib/python3.12/site-packages/dnf/base.py
usr/lib/python3.12/site-packages/dnf/callback.py
usr/lib/python3.12/site-packages/dnf/cli/
usr/lib/python3.12/site-packages/dnf/cli/__init__.py
usr/lib/python3.12/site-packages/dnf/cli/aliases.py
usr/lib/python3.12/site-packages/dnf/cli/cli.py
usr/lib/python3.12/site-packages/dnf/cli/commands/
usr/lib/python3.12/site-packages/dnf/cli/commands/__init__.py
usr/lib/python3.12/site-packages/dnf/cli/commands/alias.py
usr/lib/python3.12/site-packages/dnf/cli/commands/autoremove.py
usr/lib/python3.12/site-packages/dnf/cli/commands/check.py
usr/lib/python3.12/site-packages/dnf/cli/commands/clean.py
usr/lib/python3.12/site-packages/dnf/cli/commands/deplist.py
usr/lib/python3.12/site-packages/dnf/cli/commands/distrosync.py
usr/lib/python3.12/site-packages/dnf/cli/commands/downgrade.py
usr/lib/python3.12/site-packages/dnf/cli/commands/group.py
usr/lib/python3.12/site-packages/dnf/cli/commands/history.py
usr/lib/python3.12/site-packages/dnf/cli/commands/install.py
usr/lib/python3.12/site-packages/dnf/cli/commands/makecache.py
usr/lib/python3.12/site-packages/dnf/cli/commands/mark.py
usr/lib/python3.12/site-packages/dnf/cli/commands/module.py
usr/lib/python3.12/site-packages/dnf/cli/commands/reinstall.py
usr/lib/python3.12/site-packages/dnf/cli/commands/remove.py
usr/lib/python3.12/site-packages/dnf/cli/commands/repolist.py
usr/lib/python3.12/site-packages/dnf/cli/commands/repoquery.py
usr/lib/python3.12/site-packages/dnf/cli/commands/search.py
usr/lib/python3.12/site-packages/dnf/cli/commands/shell.py
usr/lib/python3.12/site-packages/dnf/cli/commands/swap.py
usr/lib/python3.12/site-packages/dnf/cli/commands/updateinfo.py
usr/lib/python3.12/site-packages/dnf/cli/commands/upgrade.py
usr/lib/python3.12/site-packages/dnf/cli/commands/upgrademinimal.py
usr/lib/python3.12/site-packages/dnf/cli/completion_helper.py
usr/lib/python3.12/site-packages/dnf/cli/demand.py
usr/lib/python3.12/site-packages/dnf/cli/format.py
usr/lib/python3.12/site-packages/dnf/cli/main.py
usr/lib/python3.12/site-packages/dnf/cli/option_parser.py
usr/lib/python3.12/site-packages/dnf/cli/output.py
usr/lib/python3.12/site-packages/dnf/cli/progress.py
usr/lib/python3.12/site-packages/dnf/cli/term.py
usr/lib/python3.12/site-packages/dnf/cli/utils.py
usr/lib/python3.12/site-packages/dnf/comps.py
usr/lib/python3.12/site-packages/dnf/conf/
usr/lib/python3.12/site-packages/dnf/conf/__init__.py
usr/lib/python3.12/site-packages/dnf/conf/config.py
usr/lib/python3.12/site-packages/dnf/conf/read.py
usr/lib/python3.12/site-packages/dnf/conf/substitutions.py
usr/lib/python3.12/site-packages/dnf/const.py
usr/lib/python3.12/site-packages/dnf/crypto.py
usr/lib/python3.12/site-packages/dnf/db/
usr/lib/python3.12/site-packages/dnf/db/__init__.py
usr/lib/python3.12/site-packages/dnf/db/group.py
usr/lib/python3.12/site-packages/dnf/db/history.py
usr/lib/python3.12/site-packages/dnf/dnssec.py
usr/lib/python3.12/site-packages/dnf/drpm.py
usr/lib/python3.12/site-packages/dnf/exceptions.py
usr/lib/python3.12/site-packages/dnf/goal.py
usr/lib/python3.12/site-packages/dnf/history.py
usr/lib/python3.12/site-packages/dnf/i18n.py
usr/lib/python3.12/site-packages/dnf/lock.py
usr/lib/python3.12/site-packages/dnf/logging.py
usr/lib/python3.12/site-packages/dnf/match_counter.py
usr/lib/python3.12/site-packages/dnf/module/
usr/lib/python3.12/site-packages/dnf/module/__init__.py
usr/lib/python3.12/site-packages/dnf/module/exceptions.py
usr/lib/python3.12/site-packages/dnf/module/module_base.py
usr/lib/python3.12/site-packages/dnf/package.py
usr/lib/python3.12/site-packages/dnf/persistor.py
usr/lib/python3.12/site-packages/dnf/plugin.py
usr/lib/python3.12/site-packages/dnf/pycomp.py
usr/lib/python3.12/site-packages/dnf/query.py
usr/lib/python3.12/site-packages/dnf/repo.py
usr/lib/python3.12/site-packages/dnf/repodict.py
usr/lib/python3.12/site-packages/dnf/rpm/
usr/lib/python3.12/site-packages/dnf/rpm/__init__.py
usr/lib/python3.12/site-packages/dnf/rpm/connection.py
usr/lib/python3.12/site-packages/dnf/rpm/error.py
usr/lib/python3.12/site-packages/dnf/rpm/miscutils.py
usr/lib/python3.12/site-packages/dnf/rpm/transaction.py
usr/lib/python3.12/site-packages/dnf/sack.py
usr/lib/python3.12/site-packages/dnf/selector.py
usr/lib/python3.12/site-packages/dnf/subject.py
usr/lib/python3.12/site-packages/dnf/transaction.py
usr/lib/python3.12/site-packages/dnf/transaction_sr.py
usr/lib/python3.12/site-packages/dnf/util.py
usr/lib/python3.12/site-packages/dnf/yum/
usr/lib/python3.12/site-packages/dnf/yum/__init__.py
usr/lib/python3.12/site-packages/dnf/yum/misc.py
usr/lib/python3.12/site-packages/dnf/yum/rpmtrans.py
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dnf-automatic-download.service
usr/lib/systemd/system/dnf-automatic-download.timer
usr/lib/systemd/system/dnf-automatic-install.service
usr/lib/systemd/system/dnf-automatic-install.timer
usr/lib/systemd/system/dnf-automatic-notifyonly.service
usr/lib/systemd/system/dnf-automatic-notifyonly.timer
usr/lib/systemd/system/dnf-automatic.service
usr/lib/systemd/system/dnf-automatic.timer
usr/lib/systemd/system/dnf-makecache.service
usr/lib/systemd/system/dnf-makecache.timer
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/dnf.conf
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/dnf-3
usr/share/doc/
usr/share/doc/dnf/
usr/share/doc/dnf/README.rst
usr/share/locale/
usr/share/locale/ar/
usr/share/locale/ar/LC_MESSAGES/
usr/share/locale/ar/LC_MESSAGES/dnf.mo
usr/share/locale/be/
usr/share/locale/be/LC_MESSAGES/
usr/share/locale/be/LC_MESSAGES/dnf.mo
usr/share/locale/bg/
usr/share/locale/bg/LC_MESSAGES/
usr/share/locale/bg/LC_MESSAGES/dnf.mo
usr/share/locale/bn/
usr/share/locale/bn/LC_MESSAGES/
usr/share/locale/bn/LC_MESSAGES/dnf.mo
usr/share/locale/bn_IN/
usr/share/locale/bn_IN/LC_MESSAGES/
usr/share/locale/bn_IN/LC_MESSAGES/dnf.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/dnf.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/dnf.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/dnf.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/dnf.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/dnf.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/dnf.mo
usr/share/locale/eo/
usr/share/locale/eo/LC_MESSAGES/
usr/share/locale/eo/LC_MESSAGES/dnf.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/dnf.mo
usr/share/locale/eu/
usr/share/locale/eu/LC_MESSAGES/
usr/share/locale/eu/LC_MESSAGES/dnf.mo
usr/share/locale/fa/
usr/share/locale/fa/LC_MESSAGES/
usr/share/locale/fa/LC_MESSAGES/dnf.mo
usr/share/locale/fi/
usr/share/locale/fi/LC_MESSAGES/
usr/share/locale/fi/LC_MESSAGES/dnf.mo
usr/share/locale/fil/
usr/share/locale/fil/LC_MESSAGES/
usr/share/locale/fil/LC_MESSAGES/dnf.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/dnf.mo
usr/share/locale/fur/
usr/share/locale/fur/LC_MESSAGES/
usr/share/locale/fur/LC_MESSAGES/dnf.mo
usr/share/locale/gd/
usr/share/locale/gd/LC_MESSAGES/
usr/share/locale/gd/LC_MESSAGES/dnf.mo
usr/share/locale/gu/
usr/share/locale/gu/LC_MESSAGES/
usr/share/locale/gu/LC_MESSAGES/dnf.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/dnf.mo
usr/share/locale/hi/
usr/share/locale/hi/LC_MESSAGES/
usr/share/locale/hi/LC_MESSAGES/dnf.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/dnf.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/dnf.mo
usr/share/locale/id/
usr/share/locale/id/LC_MESSAGES/
usr/share/locale/id/LC_MESSAGES/dnf.mo
usr/share/locale/ie/
usr/share/locale/ie/LC_MESSAGES/
usr/share/locale/ie/LC_MESSAGES/dnf.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/dnf.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/dnf.mo
usr/share/locale/ka/
usr/share/locale/ka/LC_MESSAGES/
usr/share/locale/ka/LC_MESSAGES/dnf.mo
usr/share/locale/kk/
usr/share/locale/kk/LC_MESSAGES/
usr/share/locale/kk/LC_MESSAGES/dnf.mo
usr/share/locale/ko/
usr/share/locale/ko/LC_MESSAGES/
usr/share/locale/ko/LC_MESSAGES/dnf.mo
usr/share/locale/lt/
usr/share/locale/lt/LC_MESSAGES/
usr/share/locale/lt/LC_MESSAGES/dnf.mo
usr/share/locale/ml/
usr/share/locale/ml/LC_MESSAGES/
usr/share/locale/ml/LC_MESSAGES/dnf.mo
usr/share/locale/mr/
usr/share/locale/mr/LC_MESSAGES/
usr/share/locale/mr/LC_MESSAGES/dnf.mo
usr/share/locale/ms/
usr/share/locale/ms/LC_MESSAGES/
usr/share/locale/ms/LC_MESSAGES/dnf.mo
usr/share/locale/nb/
usr/share/locale/nb/LC_MESSAGES/
usr/share/locale/nb/LC_MESSAGES/dnf.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/dnf.mo
usr/share/locale/om/
usr/share/locale/om/LC_MESSAGES/
usr/share/locale/om/LC_MESSAGES/dnf.mo
usr/share/locale/or/
usr/share/locale/or/LC_MESSAGES/
usr/share/locale/or/LC_MESSAGES/dnf.mo
usr/share/locale/pa/
usr/share/locale/pa/LC_MESSAGES/
usr/share/locale/pa/LC_MESSAGES/dnf.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/dnf.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/dnf.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/dnf.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/dnf.mo
usr/share/locale/si/
usr/share/locale/si/LC_MESSAGES/
usr/share/locale/si/LC_MESSAGES/dnf.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/dnf.mo
usr/share/locale/sl/
usr/share/locale/sl/LC_MESSAGES/
usr/share/locale/sl/LC_MESSAGES/dnf.mo
usr/share/locale/sq/
usr/share/locale/sq/LC_MESSAGES/
usr/share/locale/sq/LC_MESSAGES/dnf.mo
usr/share/locale/sr/
usr/share/locale/sr/LC_MESSAGES/
usr/share/locale/sr/LC_MESSAGES/dnf.mo
usr/share/locale/sr@latin/
usr/share/locale/sr@latin/LC_MESSAGES/
usr/share/locale/sr@latin/LC_MESSAGES/dnf.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/dnf.mo
usr/share/locale/th/
usr/share/locale/th/LC_MESSAGES/
usr/share/locale/th/LC_MESSAGES/dnf.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/dnf.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/dnf.mo
usr/share/locale/ur/
usr/share/locale/ur/LC_MESSAGES/
usr/share/locale/ur/LC_MESSAGES/dnf.mo
usr/share/locale/vi/
usr/share/locale/vi/LC_MESSAGES/
usr/share/locale/vi/LC_MESSAGES/dnf.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/dnf.mo
usr/share/locale/zh_Hans/
usr/share/locale/zh_Hans/LC_MESSAGES/
usr/share/locale/zh_Hans/LC_MESSAGES/dnf.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/dnf.mo
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/dnf4-transaction-json.5.gz
usr/share/man/man5/dnf4.conf.5.gz
usr/share/man/man7/
usr/share/man/man7/dnf4.modularity.7.gz
usr/share/man/man8/
usr/share/man/man8/dnf-automatic.8.gz
usr/share/man/man8/dnf4.8.gz
usr/share/man/man8/yum2dnf.8.gz

%BACKUP%
etc/dnf/automatic.conf	2ad441cf6c5b00182d1621bafd89f4bc
etc/dnf/dnf.conf	29e9686297eae56f9540fb9b56ff6ec2
etc/dnf/aliases.d/zypper.conf	f4c1cebe546d13742bfa0d33ccb7424e
etc/dnf/protected.d/dnf.conf	(null)
etc/dnf/protected.d/yum.conf	023e84808d70a20ab9211394aed061d2
etc/libreport/events.d/collect_dnf.conf	(null)
etc/logrotate.d/dnf	9fe60ffbffbc1c5504428867ad4ec255

