# shell/systemui-extra/liblingmosdk-system/src/logrotate/logrotate.cron
30 9 * * * /usr/sbin/logrotate -f /etc/logrotate.d/kylog-rotate/*

