# shell/frameworks/kglobalacceld/src/lingmo-kglobalaccel.service.in
[Unit]
Description=KDE Global Shortcuts Server
PartOf=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/kglobalacceld
BusName=org.kde.kglobalaccel
Slice=background.slice
TimeoutSec=5sec
Restart=on-failure
# Exit status 1 is used when the X11 connection drops
RestartPreventExitStatus=1
