# nxshell/qtbased/deepin-service-manager/misc/deepin-service-manager.service.in
[Unit]
Description=deepin-service-manager
After=dbus.service

[Service]
ExecStart=/usr/bin/deepin-service-manager
Restart=always
RestartSec=3
@SYSTEMD_SLICE@

ReadOnlyPaths=/usr/share/deepin-service-manager/
ReadOnlyPaths=/usr/lib/deepin-service-manager/

DevicePolicy=closed

ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
PrivateNetwork=yes
ProtectHostname=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes

@SYSTEMD_INSTALL@
