# shell/frameworks/kded/src/plasma-kded6.service.in
[Unit]
Description=KDE Daemon 6
PartOf=graphical-session.target
After=plasma-kcminit.service

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/kded6
BusName=org.kde.kded6
Slice=session.slice
TimeoutStopSec=5sec
