# shell/frameworks/kactivitymanagerd/src/lingmo-kactivitymanagerd.service.in
[Unit]
Description=KActivityManager Activity manager Service
PartOf=graphical-session.target
Before=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/kactivitymanagerd
BusName=org.kde.ActivityManager
Slice=background.slice
