# init/data_base/upower-1.90.4-1/files
%FILES%
etc/
etc/UPower/
etc/UPower/UPower.conf
usr/
usr/bin/
usr/bin/upower
usr/include/
usr/include/libupower-glib/
usr/include/libupower-glib/up-autocleanups.h
usr/include/libupower-glib/up-client.h
usr/include/libupower-glib/up-device.h
usr/include/libupower-glib/up-history-item.h
usr/include/libupower-glib/up-stats-item.h
usr/include/libupower-glib/up-types.h
usr/include/libupower-glib/up-version.h
usr/include/libupower-glib/upower.h
usr/lib/
usr/lib/girepository-1.0/
usr/lib/girepository-1.0/UPowerGlib-1.0.typelib
usr/lib/libupower-glib.so
usr/lib/libupower-glib.so.3
usr/lib/libupower-glib.so.3.1.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/upower-glib.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/upower.service
usr/lib/udev/
usr/lib/udev/hwdb.d/
usr/lib/udev/hwdb.d/95-upower-hid.hwdb
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/95-upower-hid.rules
usr/lib/udev/rules.d/95-upower-wup.rules
usr/lib/upowerd
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/interfaces/
usr/share/dbus-1/interfaces/org.freedesktop.UPower.Device.xml
usr/share/dbus-1/interfaces/org.freedesktop.UPower.KbdBacklight.xml
usr/share/dbus-1/interfaces/org.freedesktop.UPower.xml
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.freedesktop.UPower.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/org.freedesktop.UPower.conf
usr/share/gir-1.0/
usr/share/gir-1.0/UPowerGlib-1.0.gir
usr/share/gtk-doc/
usr/share/gtk-doc/html/
usr/share/gtk-doc/html/UPower/
usr/share/gtk-doc/html/UPower/Device.html
usr/share/gtk-doc/html/UPower/KbdBacklight.html
usr/share/gtk-doc/html/UPower/UPower-up-types.html
usr/share/gtk-doc/html/UPower/UPower.7.html
usr/share/gtk-doc/html/UPower/UPower.devhelp2
usr/share/gtk-doc/html/UPower/UPower.html
usr/share/gtk-doc/html/UPower/UpClient.html
usr/share/gtk-doc/html/UPower/UpDevice.html
usr/share/gtk-doc/html/UPower/UpHistoryItem.html
usr/share/gtk-doc/html/UPower/UpStatsItem.html
usr/share/gtk-doc/html/UPower/annotation-glossary.html
usr/share/gtk-doc/html/UPower/home.png
usr/share/gtk-doc/html/UPower/index.html
usr/share/gtk-doc/html/UPower/ix01.html
usr/share/gtk-doc/html/UPower/left-insensitive.png
usr/share/gtk-doc/html/UPower/left.png
usr/share/gtk-doc/html/UPower/libupower-glib-helpers.html
usr/share/gtk-doc/html/UPower/libupower-glib.html
usr/share/gtk-doc/html/UPower/license.html
usr/share/gtk-doc/html/UPower/ref-dbus.html
usr/share/gtk-doc/html/UPower/right-insensitive.png
usr/share/gtk-doc/html/UPower/right.png
usr/share/gtk-doc/html/UPower/style.css
usr/share/gtk-doc/html/UPower/tools-fileformats.html
usr/share/gtk-doc/html/UPower/up-insensitive.png
usr/share/gtk-doc/html/UPower/up.png
usr/share/gtk-doc/html/UPower/upower.1.html
usr/share/gtk-doc/html/UPower/upowerd.8.html
usr/share/locale/
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/upower.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/upower.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/upower.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/upower.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/upower.1.gz
usr/share/man/man7/
usr/share/man/man7/UPower.7.gz
usr/share/man/man8/
usr/share/man/man8/upowerd.8.gz
var/
var/lib/
var/lib/upower/

%BACKUP%
etc/UPower/UPower.conf	2b098b3a897395cf380f7b282eb97917

