# shell/systemui-extra/liblingmosdk-base/src/log/logrotate.cron
00 0 * * * /usr/bin/logrotate -s /tmp/kylog-rotate-tmp /etc/lingmosdk/lingmosdk-base/kylog-rotate-default

