# nxshell/module/lingmo-service-manager/misc/lingmo-service-plugin@.service.in
[Unit]
Description=lingmo-service-manager
After=dbus.service @EXTRA_AFTER@

[Service]
ExecStart=/usr/bin/lingmo-service-manager -n %i
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@
