# nxshell/module/lingmo-service-manager/misc/lingmo-service-manager.service.in
[Unit]
Description=lingmo-service-manager
@REQUISITE@
@PARTOF@
@BEFORE@

After=dbus.service @EXTRA_AFTER@
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart=/usr/bin/lingmo-service-manager
Restart=on-failure
RestartSec=3
@SYSTEMD_SLICE@

@SYSTEMD_INSTALL@
