# shell/frameworks/kwallet-pam/plasma-kwallet-pam.service.in
[Unit]
Description=Unlock kwallet from pam credentials
PartOf=graphical-session.target
Before=plasma-plasmashell.service plasma-kded6.service
After=plasma-kwin_wayland.service plasma-kcminit.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/pam_kwallet_init
Type=simple
Slice=background.slice
Restart=no
