# shell/systemui-extra/lingmo-media/debian/lingmo-media-server.postinst
#!/bin/sh

set -e

#DEBHELPER#

case "${1}" in
	configure)
		systemctl enable lingmo-media-server.service
		echo "start lingmo-media-server"
		;;

	abort-upgrade|abort-remove|abort-deconfigure)
		;;

	*)
		echo "postinst called with unknown argument \`${1}'" >&2
		exit 1
		;;
esac



exit 0
