# init/data_base/tpm2-abrmd-3.0.0-1/files
%FILES%
usr/
usr/bin/
usr/bin/tpm2-abrmd
usr/include/
usr/include/tss2/
usr/include/tss2/tss2-tcti-tabrmd.h
usr/lib/
usr/lib/libtss2-tcti-tabrmd.so
usr/lib/libtss2-tcti-tabrmd.so.0
usr/lib/libtss2-tcti-tabrmd.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/tss2-tcti-tabrmd.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/tpm2-abrmd.service
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/com.intel.tss2.Tabrmd.service
usr/share/dbus-1/system.d/
usr/share/dbus-1/system.d/tpm2-abrmd.conf
usr/share/licenses/
usr/share/licenses/tpm2-abrmd/
usr/share/licenses/tpm2-abrmd/LICENSE
usr/share/man/
usr/share/man/man3/
usr/share/man/man3/Tss2_Tcti_Tabrmd_Init.3.gz
usr/share/man/man7/
usr/share/man/man7/tss2-tcti-tabrmd.7.gz
usr/share/man/man8/
usr/share/man/man8/tpm2-abrmd.8.gz

