# init/data_base/hostapd-2.10-4/files
%FILES%
etc/
etc/hostapd/
etc/hostapd/hostapd.accept
etc/hostapd/hostapd.conf
etc/hostapd/hostapd.deny
etc/hostapd/hostapd.eap_user
etc/hostapd/hostapd.radius_clients
etc/hostapd/hostapd.vlan
etc/hostapd/hostapd.wpa_psk
usr/
usr/bin/
usr/bin/hostapd
usr/bin/hostapd_cli
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/hostapd.service
usr/lib/systemd/system/hostapd@.service
usr/share/
usr/share/doc/
usr/share/doc/hostapd/
usr/share/doc/hostapd/ChangeLog
usr/share/doc/hostapd/README
usr/share/doc/hostapd/README-MULTI-AP
usr/share/doc/hostapd/README-WPS
usr/share/doc/hostapd/hlr_auc_gw.milenage_db
usr/share/doc/hostapd/hlr_auc_gw.txt
usr/share/doc/hostapd/hostapd.sim_db
usr/share/doc/hostapd/wired.conf
usr/share/licenses/
usr/share/licenses/hostapd/
usr/share/licenses/hostapd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/hostapd_cli.1.gz
usr/share/man/man8/
usr/share/man/man8/hostapd.8.gz
var/
var/lib/
var/lib/hostapd/

%BACKUP%
etc/hostapd/hostapd.accept	d26c8e028ec4bb62b913ec45699252a7
etc/hostapd/hostapd.conf	fd474f24780c5458075ae0b21f04dbde
etc/hostapd/hostapd.deny	81976bea013a04fbaa0453fa54ed3f3a
etc/hostapd/hostapd.eap_user	3dbbad272b1cffaf454a3cb4db42a232
etc/hostapd/hostapd.radius_clients	13451297c5b591287bdd0217970a2900
etc/hostapd/hostapd.vlan	2a9a48c99a6fa243730a1ffb6a8a0b0f
etc/hostapd/hostapd.wpa_psk	f7c3edd4086c17a1c09b603d0144d9a3

