# nxshell/module/lingmo-diskmanager/debian/lingmo-diskmanager.postinst
#!/bin/sh
set -e
deb-systemd-helper disable 'diskmanager-daemon.service' >/dev/null || true
