# nxshell/qtbased/dde-application-wizard/systemd/org.deepin.dde.daemon.Launcher1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=Legacy DDE Application Manager launcher daemon compat background service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/dde-application-wizard-daemon-compat
Type=dbus
BusName=org.deepin.dde.daemon.Launcher1
Restart=on-failure
RestartSec=1s
