# nxshell/qtbased/dde-application-manager/apps/app-update-notifier/misc/systemd/system/org.desktopspec.ApplicationUpdateNotifier1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=Notify by signal when Application changed by package manager.

[Service]
Type=dbus
User=root
Group=root
BusName=org.desktopspec.ApplicationUpdateNotifier1
ExecStart=@CMAKE_INSTALL_FULL_LIBEXECDIR@/deepin/application-manager/app-update-notifier
