# shell/systemui-extra/lingmo-process-manager/daemon/configs/lingmo-process-manager-daemon.service
[Unit]
Description=lingmo process manager daemon
After=display-manager.service

[Service]
Type=dbus
BusName=com.lingmo.ProcessManagerDaemon
ExecStart=/usr/bin/lingmo-process-manager-daemon

[Install]
WantedBy=multi-user.target