# nxshell/qtbased/dde-display/misc/systemd/dde-display.service
[Unit]
Description=Deepin Display service

[Service]
Type=DBus
BusName=org.deepin.dde.Display1
ExecStart=/usr/bin/dde-display
