# nxshell/module/ocean-control-center/misc/systemd/ocean-control-center.service.in
[Unit]
Description=OCEAN Control Center background service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ocean-control-center -d
Type=dbus
BusName=org.lingmo.ocean.ControlCenter1
Slice=app.slice
