# build/lingmo-config/common/includes.chroot_after_packages/usr/lib/systemd/system/lingmo-watchdog.service
[Unit]
Description=Lingmo Desktop Watchdog
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/lingmo-watchdog
Restart=always

[Install]
WantedBy=graphical.target 
