# nxshell/qtbased/dde-appearance/misc/systemd/dde-appearance.service.in
[Unit]
Description=Deepin Appearance service

Requisite=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target

[Service]
Type=dbus
BusName=org.deepin.dde.Appearance1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-appearance
