# nxshell/module/ocean-clipboard/misc/ocean-clipboard-daemon.service.in
[Unit]
Description=ocean-clipboard-daemon

[Service]
Type=dbus
BusName=org.lingmo.ocean.daemon.Clipboard1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ocean-clipboard-daemon
Slice=session.slice
