# init/data_base/nfs-utils-2.6.4-1/files
%FILES%
etc/
etc/exports
etc/exports.d/
etc/nfs.conf
etc/nfsmount.conf
etc/request-key.d/
etc/request-key.d/id_resolver.conf
usr/
usr/bin/
usr/bin/blkmapd
usr/bin/exportfs
usr/bin/fsidd
usr/bin/mount.nfs
usr/bin/mount.nfs4
usr/bin/mountstats
usr/bin/nfsconf
usr/bin/nfsdcld
usr/bin/nfsdclddb
usr/bin/nfsdclnts
usr/bin/nfsdcltrack
usr/bin/nfsidmap
usr/bin/nfsiostat
usr/bin/nfsrahead
usr/bin/nfsstat
usr/bin/nfsv4.exportd
usr/bin/rpc.gssd
usr/bin/rpc.idmapd
usr/bin/rpc.mountd
usr/bin/rpc.nfsd
usr/bin/rpc.statd
usr/bin/rpcctl
usr/bin/rpcdebug
usr/bin/showmount
usr/bin/sm-notify
usr/bin/start-statd
usr/bin/umount.nfs
usr/bin/umount.nfs4
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-generators/
usr/lib/systemd/system-generators/nfs-server-generator
usr/lib/systemd/system-generators/rpc-pipefs-generator
usr/lib/systemd/system/
usr/lib/systemd/system/auth-rpcgss-module.service
usr/lib/systemd/system/fsidd.service
usr/lib/systemd/system/nfs-blkmap.service
usr/lib/systemd/system/nfs-client.target
usr/lib/systemd/system/nfs-idmapd.service
usr/lib/systemd/system/nfs-mountd.service
usr/lib/systemd/system/nfs-server.service
usr/lib/systemd/system/nfs-utils.service
usr/lib/systemd/system/nfsdcld.service
usr/lib/systemd/system/nfsv4-exportd.service
usr/lib/systemd/system/nfsv4-server.service
usr/lib/systemd/system/proc-fs-nfsd.mount
usr/lib/systemd/system/rpc-gssd.service
usr/lib/systemd/system/rpc-statd-notify.service
usr/lib/systemd/system/rpc-statd.service
usr/lib/systemd/system/rpc_pipefs.target
usr/lib/systemd/system/var-lib-nfs-rpc_pipefs.mount
usr/lib/sysusers.d/
usr/lib/sysusers.d/rpcuser.conf
usr/lib/udev/
usr/lib/udev/rules.d/
usr/lib/udev/rules.d/60-nfs.rules
usr/lib/udev/rules.d/99-nfs.rules
usr/share/
usr/share/doc/
usr/share/doc/nfs-utils/
usr/share/doc/nfs-utils/NEWS
usr/share/doc/nfs-utils/README
usr/share/doc/nfs-utils/README.systemd
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/exports.5.gz
usr/share/man/man5/nfs.5.gz
usr/share/man/man5/nfs.conf.5.gz
usr/share/man/man5/nfsmount.conf.5.gz
usr/share/man/man5/nfsrahead.5.gz
usr/share/man/man7/
usr/share/man/man7/nfs.systemd.7.gz
usr/share/man/man7/nfsd.7.gz
usr/share/man/man8/
usr/share/man/man8/blkmapd.8.gz
usr/share/man/man8/exportd.8.gz
usr/share/man/man8/exportfs.8.gz
usr/share/man/man8/gssd.8.gz
usr/share/man/man8/idmapd.8.gz
usr/share/man/man8/mount.nfs.8.gz
usr/share/man/man8/mountd.8.gz
usr/share/man/man8/mountstats.8.gz
usr/share/man/man8/nfsconf.8.gz
usr/share/man/man8/nfsd.8.gz
usr/share/man/man8/nfsdcld.8.gz
usr/share/man/man8/nfsdclddb.8.gz
usr/share/man/man8/nfsdclnts.8.gz
usr/share/man/man8/nfsdcltrack.8.gz
usr/share/man/man8/nfsidmap.8.gz
usr/share/man/man8/nfsiostat.8.gz
usr/share/man/man8/nfsstat.8.gz
usr/share/man/man8/nfsv4.exportd.8.gz
usr/share/man/man8/rpc.gssd.8.gz
usr/share/man/man8/rpc.idmapd.8.gz
usr/share/man/man8/rpc.mountd.8.gz
usr/share/man/man8/rpc.nfsd.8.gz
usr/share/man/man8/rpc.sm-notify.8.gz
usr/share/man/man8/rpc.statd.8.gz
usr/share/man/man8/rpcctl.8.gz
usr/share/man/man8/rpcdebug.8.gz
usr/share/man/man8/showmount.8.gz
usr/share/man/man8/sm-notify.8.gz
usr/share/man/man8/statd.8.gz
usr/share/man/man8/umount.nfs.8.gz
var/
var/lib/
var/lib/nfs/
var/lib/nfs/etab
var/lib/nfs/rmtab
var/lib/nfs/rpc_pipefs/
var/lib/nfs/statd/
var/lib/nfs/statd/sm.bak/
var/lib/nfs/statd/sm/
var/lib/nfs/statd/state
var/lib/nfs/v4recovery/

%BACKUP%
etc/exports	cf5b3d6de892a5e9960af7b47faccc8b
etc/nfs.conf	5ad00b2a45a46de686fa44009215a205
etc/nfsmount.conf	8944afe67f93ec7e709190cdf8fe26ca
var/lib/nfs/etab	d41d8cd98f00b204e9800998ecf8427e
var/lib/nfs/rmtab	d41d8cd98f00b204e9800998ecf8427e

