# shell/ai/lingmo-ai-runtime/configs/lib/systemd/system/lingmo-ai-runtime.service
[Unit]
Description=Lingmo Ai System Runtime
After=display-manager.service

[Service]
ExecStart=/usr/sbin/lingmo-ai-runtime
Restart=always

[Install]
WantedBy=multi-user.target