# nxshell/qtbased/deepin-log-viewer/logViewerService/assets/data/deepin-log-viewer-daemon.service
[Unit]
Description=Deepin Log Viewer Daemon

[Service]
Type=dbus
BusName=com.deepin.logviewer
ExecStart=/usr/lib/deepin-daemon/log-view-service
CapabilityBoundingSet=~CAP_NET_RAW
MemoryLimit=8G

[Install]
WantedBy=multi-user.target
