# app/lingmo-keyring/wkd_sync/archlinux-keyring-wkd-sync.timer
[Unit]
Description=Refresh existing PGP keys of archlinux-keyring regularly

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1week
