# init/data_base/containerd-1.7.20-1/files
%FILES%
usr/
usr/bin/
usr/bin/containerd
usr/bin/containerd-shim
usr/bin/containerd-shim-runc-v1
usr/bin/containerd-shim-runc-v2
usr/bin/containerd-stress
usr/bin/ctr
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/containerd.service
usr/share/
usr/share/bash-completion/
usr/share/bash-completion/completions/
usr/share/bash-completion/completions/ctr
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/containerd-config.toml.5.gz
usr/share/man/man8/
usr/share/man/man8/containerd-config.8.gz
usr/share/man/man8/containerd.8.gz
usr/share/man/man8/ctr.8.gz
usr/share/zsh/
usr/share/zsh/site-functions/
usr/share/zsh/site-functions/_ctr

