# nxshell/qtbased/dde-session/tools/dde-version-checker/dde-version-checker.service
[Unit]
Description=dde-version-checker service
CollectMode=inactive-or-failed

Requisite=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target

[Service]
Type=simple
ExecStart=/usr/bin/dde-version-checker
