# nxshell/module/ocean-session/tools/ocean-login-reminder/ocean-login-reminder.service
[Unit]
Description=ocean-login-reminder service
CollectMode=inactive-or-failed

Wants=ocean-osd.target
After=ocean-osd.target

[Service]
Type=simple
ExecStart=/usr/bin/ocean-login-reminder
