# init/data_base/monit-5.34.0-1/files
%FILES%
etc/
etc/monitrc
usr/
usr/bin/
usr/bin/monit
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/monit.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/monit.1.gz

%BACKUP%
etc/monitrc	6b349a0e54c9a054bccbe803735162f5

