# nxshell/module/ocean-shell/misc/ocean-shell-plugin@.service.in
[Unit]
Description=Manage a plugin for ocean-shell.

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ocean-shell -p %I

[Install]
WantedBy=multi-user.target
