# nxshell/module/ocean-shell/misc/ocean-shell@.service.in
[Unit]
Description=Manage an category for ocean-shell.

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ocean-shell -C %I

[Install]
WantedBy=multi-user.target
