# nxshell/qtbased/deepin-log-viewer/application/configs/coredump-reporter.service
[Unit]
Description=deepin-log-viewer coredump report activities

[Service]

ExecStart=/usr/bin/deepin-log-viewer --reportcoredump
CapabilityBoundingSet=~
MemoryLimit=8G

[Install]
WantedBy=multi-user.target
