# init/data_base/dhcp-4.4.3.P1-3/files
%FILES%
etc/
etc/dhcpd.conf
etc/dhcpd6.conf
usr/
usr/bin/
usr/bin/dhcpd
usr/bin/dhcrelay
usr/bin/omshell
usr/include/
usr/include/dhcpctl/
usr/include/dhcpctl/dhcpctl.h
usr/include/omapip/
usr/include/omapip/alloc.h
usr/include/omapip/buffer.h
usr/include/omapip/convert.h
usr/include/omapip/hash.h
usr/include/omapip/isclib.h
usr/include/omapip/omapip.h
usr/include/omapip/omapip_p.h
usr/include/omapip/result.h
usr/include/omapip/trace.h
usr/lib/
usr/lib/libdhcp.a
usr/lib/libdhcpctl.a
usr/lib/libomapi.a
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dhcpd4.service
usr/lib/systemd/system/dhcpd6.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/dhcp.conf
usr/lib/tmpfiles.d/
usr/lib/tmpfiles.d/dhcp.conf
usr/share/
usr/share/licenses/
usr/share/licenses/dhcp/
usr/share/licenses/dhcp/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/omshell.1.gz
usr/share/man/man3/
usr/share/man/man3/dhcpctl.3.gz
usr/share/man/man3/omapi.3.gz
usr/share/man/man5/
usr/share/man/man5/dhcp-eval.5.gz
usr/share/man/man5/dhcp-options.5.gz
usr/share/man/man5/dhcpd.conf.5.gz
usr/share/man/man5/dhcpd.leases.5.gz
usr/share/man/man8/
usr/share/man/man8/dhcpd.8.gz
usr/share/man/man8/dhcrelay.8.gz

%BACKUP%
etc/dhcpd.conf	09b995321b09d019646b5525db046c15
etc/dhcpd6.conf	5a8f0cfcbb5975434c6b81930c239c45

