# init/data_base/seatd-0.8.0-2/files
%FILES%
usr/
usr/bin/
usr/bin/seatd
usr/bin/seatd-launch
usr/include/
usr/include/libseat.h
usr/lib/
usr/lib/libseat.so
usr/lib/libseat.so.1
usr/lib/pkgconfig/
usr/lib/pkgconfig/libseat.pc
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/seatd.service
usr/lib/sysusers.d/
usr/lib/sysusers.d/seatd.conf
usr/share/
usr/share/licenses/
usr/share/licenses/seatd/
usr/share/licenses/seatd/LICENSE
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/seatd-launch.1.gz
usr/share/man/man1/seatd.1.gz

