# init/data_base/dbus-broker-units-36-4/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/dbus.service
usr/lib/systemd/user/
usr/lib/systemd/user/dbus.service

