# core32/ksplash/ksplashqml/lingmo-ksplash.service.in
[Unit]
Description=Splash screen shown during boot
PartOf=graphical-session.target
After=lingmo-kwin_wayland.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ksplashqml
Type=oneshot
Slice=background.slice
TimeoutSec=5sec