# init/data_base/dbus-broker-36-4/files
%FILES%
usr/
usr/bin/
usr/bin/dbus-broker
usr/bin/dbus-broker-launch
usr/lib/
usr/lib/systemd/
usr/lib/systemd/catalog/
usr/lib/systemd/catalog/dbus-broker-launch.catalog
usr/lib/systemd/catalog/dbus-broker.catalog
usr/lib/systemd/system/
usr/lib/systemd/system/dbus-broker.service
usr/lib/systemd/user/
usr/lib/systemd/user/dbus-broker.service
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/dbus-broker-launch.1.gz
usr/share/man/man1/dbus-broker.1.gz

