# core32/startshell/systemd/lingmo-core.target
[Unit]
Description=Lingmo OS Workspace Core
Wants=lingmo-startshell.service graphical-session-pre.target
Requires=lingmo-ksmserver.service
After=graphical-session-pre.target lingmo-kwin_wayland.service
RefuseManualStart=yes
StopWhenUnneeded=true
