# core32/startshell/systemd/lingmo-workspace-x11.target
[Unit]
Wants=lingmo-kwin_x11.service
Requires=lingmo-workspace.target
BindsTo=lingmo-ksmserver.service
