# mainRelease/config/includes.chroot_after_packages/etc/systemd/system/first_setup.service
[Unit]
Description=Run a Custom Script at Startup
After=default.target

[Service]
ExecStart=/etc/setup/set_up.sh

[Install]
WantedBy=default.target
