# shell/systemui-extra/lingmo-media/configuration/99-lingmo-media.rules
# do not edit this file, it will be overwritten on update
# This file is part of LINGMO-MEDIA
# Author: Zoey

SUBSYSTEM!="sound", GOTO="pulseaudio_end"
ACTION!="change", GOTO="pulseaudio_end"
KERNEL!="card*", GOTO="pulseaudio_end"
SUBSYSTEMS=="usb", GOTO="pulseaudio_check_usb"

# 运行声卡适配脚本
RUN+="/usr/share/lingmo-media/scripts/soundCardAdaptation.sh"

SUBSYSTEMS=="platform", DRIVERS=="thinkpad_acpi", ENV{PULSE_IGNORE}="1"

SUBSYSTEMS=="platform", DRIVERS=="ft-hda", ENV{PULSE_PROFILE_SET}="pulse-ft-hda.conf"

SUBSYSTEMS=="platform", DRIVERS=="pmdk_dp", ENV{PULSE_PROFILE_SET}="pulse-x100-hdmi.conf"

GOTO="pulseaudio_end"


LABEL="pulseaudio_check_usb"

# JBL USB
ATTRS{idVendor}=="e5b7", ATTRS{idProduct}=="0811", ENV{PULSE_PROFILE_SET}="pulse-jbl-usb-audio.conf"

# Inspur Camera
ATTRS{idVendor}=="f00f", ATTRS{idProduct}=="0e01", ENV{PULSE_PROFILE_SET}="lingmo-inspur-usb-audio.conf"

# Lenovo MH243A
ATTRS{idVendor}=="0bda", ATTRS{idProduct}=="4f18", ENV{PULSE_PROFILE_SET}="alc4080-mh243a-usb-audio.conf"

# Lenovo M90H
ATTRS{idVendor}=="17aa", ATTRS{idProduct}=="3501", ENV{PULSE_PROFILE_SET}="lenovo-m90h-usb-audio.conf"

#LENOVO P90H
ATTRS{idVendor}=="17aa", ATTRS{idProduct}=="3507", ENV{PULSE_PROFILE_SET}="lenovo-p90h-usb-audio.conf"

# MCA500 usb audio
ATTRS{idVendor}=="20b1", ATTRS{idProduct}=="3079", ENV{PULSE_PROFILE_SET}="mca500-usb-audio.conf"

# Audio fs 1.0
ATTRS{idVendor}=="31b2", ATTRS{idProduct}=="0020", ENV{PULSE_PROFILE_SET}="fs-hygon-usb-audio.conf"

# Douyu G591
ATTRS{idVendor}=="2b0d", ATTRS{idProduct}=="0031", ENV{PULSE_VOLUME_IGNORE_DB}="1"

# C-media
ATTRS{idVendor}=="0d8c", ATTRS{idProduct}=="0012", ENV{PULSE_VOLUME_IGNORE_DB}="1"

# AEEVISION Camera
ATTRS{idVendor}=="058e", ATTRS{idProduct}=="3863", ENV{PULSE_VOLUME_IGNORE_DB}="1"

# EDIFIER M16
ATTRS{idVendor}=="0000", ATTRS{idProduct}=="0201", ENV{PULSE_VOLUME_IGNORE_DB}="1"

GOTO="pulseaudio_end"


LABEL="pulseaudio_end"
