# init/data_base/radvd-2.19-2/files
%FILES%
etc/
etc/radvd.conf
usr/
usr/bin/
usr/bin/radvd
usr/bin/radvdump
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/radvd.service
usr/share/
usr/share/licenses/
usr/share/licenses/radvd/
usr/share/licenses/radvd/COPYRIGHT
usr/share/man/
usr/share/man/man5/
usr/share/man/man5/radvd.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/radvd.8.gz
usr/share/man/man8/radvdump.8.gz

%BACKUP%
etc/radvd.conf	b160317dfc3ab51fa5cced871894d6dd

