# shell/systemui-extra/liblingmosdk-system/debian/liblingmosdk-system-dbus.postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable lingmosdk-dbus.service

