# shell/systemui-extra/liblingmosdk-system/debian/liblingmosdk-logrotate.postinst
#!/bin/bash

systemctl enable lingmosdk-logrotate.service
systemctl restart lingmosdk-logrotate.service
crontab /usr/lib/lingmosdk/lingmosdk-system/logrotate.cron