# init/data_base/atop-2.10.0-1/files
%FILES%
etc/
etc/default/
etc/default/atop
usr/
usr/bin/
usr/bin/atop
usr/bin/atop-2.10.0
usr/bin/atopacctd
usr/bin/atopcat
usr/bin/atopconvert
usr/bin/atopgpud
usr/bin/atophide
usr/bin/atopsar
usr/bin/atopsar-2.10.0
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system-sleep/
usr/lib/systemd/system-sleep/atop-pm.sh
usr/lib/systemd/system/
usr/lib/systemd/system/atop-rotate.service
usr/lib/systemd/system/atop-rotate.timer
usr/lib/systemd/system/atop.service
usr/lib/systemd/system/atopacct.service
usr/lib/systemd/system/atopgpu.service
usr/lib/systemd/system/timers.target.wants/
usr/lib/systemd/system/timers.target.wants/atop-rotate.timer
usr/share/
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/atop.1.gz
usr/share/man/man1/atopcat.1.gz
usr/share/man/man1/atopconvert.1.gz
usr/share/man/man1/atophide.1.gz
usr/share/man/man1/atopsar.1.gz
usr/share/man/man5/
usr/share/man/man5/atoprc.5.gz
usr/share/man/man8/
usr/share/man/man8/atopacctd.8.gz
usr/share/man/man8/atopgpud.8.gz
var/
var/log/
var/log/atop/

%BACKUP%
etc/default/atop	dea57a98615a07f537212468ad5b0857

