# distribution/system/bookworm/lingmo-config/common/includes.chroot_befor_packages/etc/profile.d/home-local-bin.sh
case ":${PATH}:" in
  *:"$HOME/.local/bin":*) ;;
  *) export PATH="$HOME/.local/bin:$PATH" ;;
esac
