# nxshell/qtbased/dde-clipboard/misc/dde-clipboard-daemon.service.in
[Unit]
Description=dde-clipboard-daemon

[Service]
Type=dbus
BusName=org.deepin.dde.daemon.Clipboard1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-clipboard-daemon
Slice=session.slice
