# nxshell/qtbased/dde-session/tools/dde-login-reminder/dde-login-reminder.service
[Unit]
Description=dde-login-reminder service
CollectMode=inactive-or-failed

Wants=dde-osd.target
After=dde-osd.target

[Service]
Type=simple
ExecStart=/usr/bin/dde-login-reminder
