# init/data_base/busybox-1.36.1-2/files
%FILES%
usr/
usr/bin/
usr/bin/busybox
usr/share/
usr/share/doc/
usr/share/doc/busybox/
usr/share/doc/busybox/BusyBox.html
usr/share/doc/busybox/BusyBox.txt
usr/share/doc/busybox/examples/
usr/share/doc/busybox/examples/android-build
usr/share/doc/busybox/examples/bootfloppy/
usr/share/doc/busybox/examples/bootfloppy/bootfloppy.txt
usr/share/doc/busybox/examples/bootfloppy/display.txt
usr/share/doc/busybox/examples/bootfloppy/etc/
usr/share/doc/busybox/examples/bootfloppy/etc/fstab
usr/share/doc/busybox/examples/bootfloppy/etc/init.d/
usr/share/doc/busybox/examples/bootfloppy/etc/init.d/rcS
usr/share/doc/busybox/examples/bootfloppy/etc/inittab
usr/share/doc/busybox/examples/bootfloppy/etc/profile
usr/share/doc/busybox/examples/bootfloppy/mkdevs.sh
usr/share/doc/busybox/examples/bootfloppy/mkrootfs.sh
usr/share/doc/busybox/examples/bootfloppy/mksyslinux.sh
usr/share/doc/busybox/examples/bootfloppy/quickstart.txt
usr/share/doc/busybox/examples/bootfloppy/syslinux.cfg
usr/share/doc/busybox/examples/busybox.spec
usr/share/doc/busybox/examples/depmod
usr/share/doc/busybox/examples/depmod.pl
usr/share/doc/busybox/examples/devfsd.conf
usr/share/doc/busybox/examples/dnsd.conf
usr/share/doc/busybox/examples/inetd.conf
usr/share/doc/busybox/examples/inittab
usr/share/doc/busybox/examples/linux-2.6.30_proc_self_exe.patch
usr/share/doc/busybox/examples/mdev.conf
usr/share/doc/busybox/examples/mdev.conf.change_blockdev.sh
usr/share/doc/busybox/examples/mdev_fat.conf
usr/share/doc/busybox/examples/mk2knr.pl
usr/share/doc/busybox/examples/shutdown-1.0/
usr/share/doc/busybox/examples/shutdown-1.0/README
usr/share/doc/busybox/examples/shutdown-1.0/script/
usr/share/doc/busybox/examples/shutdown-1.0/script/do_shutdown
usr/share/doc/busybox/examples/shutdown-1.0/script/hardshutdown.c
usr/share/doc/busybox/examples/shutdown-1.0/script/hardshutdown.make.sh
usr/share/doc/busybox/examples/shutdown-1.0/script/shutdown
usr/share/doc/busybox/examples/shutdown-1.0/script/stop_storage
usr/share/doc/busybox/examples/shutdown-1.0/script/stop_tasks
usr/share/doc/busybox/examples/udhcp/
usr/share/doc/busybox/examples/udhcp/sample.bound
usr/share/doc/busybox/examples/udhcp/sample.deconfig
usr/share/doc/busybox/examples/udhcp/sample.nak
usr/share/doc/busybox/examples/udhcp/sample.renew
usr/share/doc/busybox/examples/udhcp/sample.script
usr/share/doc/busybox/examples/udhcp/simple.script
usr/share/doc/busybox/examples/udhcp/udhcpd.conf
usr/share/doc/busybox/examples/undeb
usr/share/doc/busybox/examples/unrpm
usr/share/doc/busybox/examples/var_service/
usr/share/doc/busybox/examples/var_service/README
usr/share/doc/busybox/examples/var_service/README_distro_proposal.txt
usr/share/doc/busybox/examples/var_service/dhcp_if/
usr/share/doc/busybox/examples/var_service/dhcp_if/convert2ipconf
usr/share/doc/busybox/examples/var_service/dhcp_if/convert2ntpconf
usr/share/doc/busybox/examples/var_service/dhcp_if/dhcp_handler
usr/share/doc/busybox/examples/var_service/dhcp_if/finish
usr/share/doc/busybox/examples/var_service/dhcp_if/log/
usr/share/doc/busybox/examples/var_service/dhcp_if/log/run
usr/share/doc/busybox/examples/var_service/dhcp_if/run
usr/share/doc/busybox/examples/var_service/dhcp_if_pinger/
usr/share/doc/busybox/examples/var_service/dhcp_if_pinger/run
usr/share/doc/busybox/examples/var_service/dhcpd_if/
usr/share/doc/busybox/examples/var_service/dhcpd_if/log/
usr/share/doc/busybox/examples/var_service/dhcpd_if/log/run
usr/share/doc/busybox/examples/var_service/dhcpd_if/run
usr/share/doc/busybox/examples/var_service/dhcpd_if/udhcpd.conf
usr/share/doc/busybox/examples/var_service/dhcpd_if/w_dumpleases
usr/share/doc/busybox/examples/var_service/dhcpd_if/w_dumpleases_countdown
usr/share/doc/busybox/examples/var_service/dnsmasq/
usr/share/doc/busybox/examples/var_service/dnsmasq/README
usr/share/doc/busybox/examples/var_service/dnsmasq/dnsmasq.conf
usr/share/doc/busybox/examples/var_service/dnsmasq/log/
usr/share/doc/busybox/examples/var_service/dnsmasq/log/run
usr/share/doc/busybox/examples/var_service/dnsmasq/run
usr/share/doc/busybox/examples/var_service/ftpd/
usr/share/doc/busybox/examples/var_service/ftpd/log/
usr/share/doc/busybox/examples/var_service/ftpd/log/run
usr/share/doc/busybox/examples/var_service/ftpd/run
usr/share/doc/busybox/examples/var_service/fw/
usr/share/doc/busybox/examples/var_service/fw/conf/
usr/share/doc/busybox/examples/var_service/fw/conf/11.22.33.44.ipconf--
usr/share/doc/busybox/examples/var_service/fw/conf/192.168.0.1.ipconf
usr/share/doc/busybox/examples/var_service/fw/conf/lo.ipconf
usr/share/doc/busybox/examples/var_service/fw/etc/
usr/share/doc/busybox/examples/var_service/fw/etc/dnsmasq_servers.conf
usr/share/doc/busybox/examples/var_service/fw/etc/hosts
usr/share/doc/busybox/examples/var_service/fw/etc/resolv.conf
usr/share/doc/busybox/examples/var_service/fw/run
usr/share/doc/busybox/examples/var_service/fw/stat
usr/share/doc/busybox/examples/var_service/getty_tty1/
usr/share/doc/busybox/examples/var_service/getty_tty1/alt08x16+unimap.fnt
usr/share/doc/busybox/examples/var_service/getty_tty1/cfg
usr/share/doc/busybox/examples/var_service/getty_tty1/koi8r_to_uni.trans
usr/share/doc/busybox/examples/var_service/getty_tty1/login.sh
usr/share/doc/busybox/examples/var_service/getty_tty1/run
usr/share/doc/busybox/examples/var_service/getty_tty1/unicode_cyrillic.keymap
usr/share/doc/busybox/examples/var_service/gpm/
usr/share/doc/busybox/examples/var_service/gpm/run
usr/share/doc/busybox/examples/var_service/httpd/
usr/share/doc/busybox/examples/var_service/httpd/log/
usr/share/doc/busybox/examples/var_service/httpd/log/run
usr/share/doc/busybox/examples/var_service/httpd/run
usr/share/doc/busybox/examples/var_service/ifplugd_if/
usr/share/doc/busybox/examples/var_service/ifplugd_if/ifplugd_handler
usr/share/doc/busybox/examples/var_service/ifplugd_if/log/
usr/share/doc/busybox/examples/var_service/ifplugd_if/log/run
usr/share/doc/busybox/examples/var_service/ifplugd_if/run
usr/share/doc/busybox/examples/var_service/inetd/
usr/share/doc/busybox/examples/var_service/inetd/inetd.conf
usr/share/doc/busybox/examples/var_service/inetd/log/
usr/share/doc/busybox/examples/var_service/inetd/log/run
usr/share/doc/busybox/examples/var_service/inetd/run
usr/share/doc/busybox/examples/var_service/nmeter/
usr/share/doc/busybox/examples/var_service/nmeter/run
usr/share/doc/busybox/examples/var_service/ntpd/
usr/share/doc/busybox/examples/var_service/ntpd/log/
usr/share/doc/busybox/examples/var_service/ntpd/log/run
usr/share/doc/busybox/examples/var_service/ntpd/ntp.script
usr/share/doc/busybox/examples/var_service/ntpd/p_log_important
usr/share/doc/busybox/examples/var_service/ntpd/run
usr/share/doc/busybox/examples/var_service/std_service_logger
usr/share/doc/busybox/examples/var_service/supplicant_if/
usr/share/doc/busybox/examples/var_service/supplicant_if/log/
usr/share/doc/busybox/examples/var_service/supplicant_if/log/run
usr/share/doc/busybox/examples/var_service/supplicant_if/run
usr/share/doc/busybox/examples/var_service/supplicant_if/wpa_supplicant.conf
usr/share/doc/busybox/examples/var_service/sview
usr/share/doc/busybox/examples/var_service/svpage
usr/share/doc/busybox/examples/var_service/tftpd/
usr/share/doc/busybox/examples/var_service/tftpd/log/
usr/share/doc/busybox/examples/var_service/tftpd/log/run
usr/share/doc/busybox/examples/var_service/tftpd/run
usr/share/doc/busybox/examples/var_service/zcip_if/
usr/share/doc/busybox/examples/var_service/zcip_if/convert2ipconf
usr/share/doc/busybox/examples/var_service/zcip_if/finish
usr/share/doc/busybox/examples/var_service/zcip_if/log/
usr/share/doc/busybox/examples/var_service/zcip_if/log/run
usr/share/doc/busybox/examples/var_service/zcip_if/run
usr/share/doc/busybox/examples/var_service/zcip_if/zcip_handler
usr/share/doc/busybox/examples/zcip.script
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/busybox.1.gz

