# shell/frameworks/lingmo-workspace/startkde/systemd/lingmo-workspace.target
[Unit]
Description=KDE Lingmo Workspace
Requires=lingmo-core.target graphical-session.target
Wants=lingmo-restoresession.service lingmo-xembedsniproxy.service lingmo-gmenudbusmenuproxy.service lingmo-powerdevil.service lingmo-ksplash-ready.service lingmo-polkit-agent.service kde-baloo.service lingmo-foreground-booster.service lingmo-kwallet-pam.service lingmo-kglobalaccel.service lingmo-kaccess.service xdg-desktop-autostart.target
BindsTo=graphical-session.target
Before=graphical-session.target xdg-desktop-autostart.target lingmo-ksplash-ready.service lingmo-restoresession.service
RefuseManualStart=yes
StopWhenUnneeded=true
