# shell/frameworks/lingmo-workspace/startkde/systemd/lingmo-workspace-wayland.target
[Unit]
Requires=lingmo-workspace.target
Requires=lingmo-kwin_wayland.service
BindsTo=lingmo-kwin_wayland.service
