# shell/frameworks/lingmo-workspace/startkde/systemd/lingmo-core.target
[Unit]
Description=KDE Lingmo Workspace Core
Wants=lingmo-lingmoshell.service lingmo-kcminit.service lingmo-kded6.service lingmo-kcminit-phase1.service graphical-session-pre.target
Requires=lingmo-ksmserver.service
After=graphical-session-pre.target lingmo-kwin_wayland.service
RefuseManualStart=yes
StopWhenUnneeded=true
