# nxshell/module/ocean-calendar/calendar-service/assets/data/ocean-calendar.service.in
[Unit]
Description = ocean-calendar date service

[Service]
ExecStart=@EXEC_OCEAN_CALENDAR@
Type=dbus
BusName=com.lingmo.Calendar
Slice=app.slice