# shell/systemui-extra/liblingmosdk-base/src/conf2/service/conf2-update-db.service
[Unit]
Description=LINGMO CONF2 UPDATE DB
DefaultDependencies=no
Before=dbus.service

[Service]
Type=oneshot
ExecStart=/usr/bin/conf2-update.py
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
