# shell/systemui-extra/liblingmosdk-base/src/conf2/service/lingmosdk-conf2.service
[Unit]
Description=LINGMO CONF2 DBUS
Before=systemd-logind.service com.lingmo.lingmosdk.SyncConfig.service

[Service]
Type=dbus
Restart=always
RestartSec=1
BusName=com.lingmo.lingmosdk.conf2
ExecStart=/usr/bin/conf2-server.py

[Install]
WantedBy=multi-user.target
